# hash256-miner.service
[Unit]
Description=HASH256 GPU Miner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/hash256-miner
EnvironmentFile=/root/hash256-miner/.env
ExecStart=/usr/bin/python3 /root/hash256-miner/miner.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
